# templates/systemd/nginx.service.d/kamailio.conf
[Unit]
Wants=kamailio.service
After=kamailio.service
