# templates/systemd/kamailio.service.d/database.conf
[Unit]
Requires=mariadb.service
Requires=rtpengine-daemon.service
After=mariadb.service rtpengine-daemon.service
